# packages/rhel/SOURCES/zabbix-agent-dbmon.service
[Unit]
Description=Zabbix Agent for Database monitoring
After=syslog.target
After=network.target

[Service]
Environment="CONFFILE=/etc/zabbix/zabbix_agentd_dbmon.conf"
EnvironmentFile=-/etc/sysconfig/zabbix-agent-dbmon
Type=forking
Restart=on-failure
PIDFile=/run/zabbix/zabbix_agentd_dbmon.pid
KillMode=control-group
ExecStart=/usr/sbin/zabbix_agentd_dbmon -c $CONFFILE
ExecStop=/bin/kill -SIGTERM $MAINPID
RestartSec=10s
User=zabbix
Group=zabbix

[Install]
WantedBy=multi-user.target
